# root/etc/init.d/eventcenter
#!/bin/sh /etc/rc.common

START=99
STOP=10

WATCHER_PID_FILE="/tmp/eventcenter_watcher.pid"

start_service() {
    local enable
    enable=$(uci -q get eventcenter.global.enable)
    [ "$enable" != "1" ] && return

    # --- OpenClash cron job setup ---
    local interval
    interval=$(uci -q get eventcenter.@monitor[0].interval 2>/dev/null)
    [ -z "$interval" ] && interval=5

    local cron_expr
    if [ "$interval" -ge 1440 ] 2>/dev/null && [ $(( interval % 1440 )) -eq 0 ] 2>/dev/null; then
        local days=$(( interval / 1440 ))
        cron_expr="0 0 */${days} * *"
    elif [ "$interval" -ge 60 ] 2>/dev/null && [ $(( interval % 60 )) -eq 0 ] 2>/dev/null; then
        local hours=$(( interval / 60 ))
        cron_expr="0 */${hours} * * *"
    else
        cron_expr="*/${interval} * * * *"
    fi

    # Remove old eventcenter cron entries and add new ones
    grep -v 'eventcenter_check' /etc/crontabs/root > /tmp/ec_cr_$$ 2>/dev/null
    echo "${cron_expr} /usr/bin/eventcenter check openclash" >> /tmp/ec_cr_$$
    cp /tmp/ec_cr_$$ /etc/crontabs/root
    rm -f /tmp/ec_cr_$$
    logger -t eventcenter "openclash cron: ${cron_expr} (interval=${interval}min)"

    # --- Node health monitoring ---
    local health_enable
    health_enable=$(uci -q get eventcenter.@health[0].enable 2>/dev/null)

    # Always remove old health cron first
    grep -v 'eventcenter_check_node_health' /etc/crontabs/root > /tmp/ec_crh_$$ 2>/dev/null
    cp /tmp/ec_crh_$$ /etc/crontabs/root
    rm -f /tmp/ec_crh_$$

    if [ "$health_enable" = "1" ]; then
        local health_interval
        health_interval=$(uci -q get eventcenter.@health[0].interval 2>/dev/null || echo "5")
        [ -z "$health_interval" ] && health_interval=5

        local health_cron
        if [ "$health_interval" -ge 1440 ] 2>/dev/null && [ $(( health_interval % 1440 )) -eq 0 ] 2>/dev/null; then
            local health_days=$(( health_interval / 1440 ))
            health_cron="0 0 */${health_days} * *"
        elif [ "$health_interval" -ge 60 ] 2>/dev/null && [ $(( health_interval % 60 )) -eq 0 ] 2>/dev/null; then
            local health_hours=$(( health_interval / 60 ))
            health_cron="0 */${health_hours} * * *"
        else
            health_cron="*/${health_interval} * * * *"
        fi

        echo "${health_cron} /usr/bin/eventcenter check node-health" >> /etc/crontabs/root
        logger -t eventcenter "node-health cron: ${health_cron} (interval=${health_interval}min)"
    fi

    # --- Watcher setup ---
    local realtime
    realtime=$(uci -q get eventcenter.@monitor[0].realtime 2>/dev/null)
    if [ "$realtime" = "1" ]; then
        # Stop existing watcher first
        if [ -f "$WATCHER_PID_FILE" ]; then
            local _old_pid
            _old_pid=$(cat "$WATCHER_PID_FILE")
            kill "$_old_pid" 2>/dev/null
            rm -f "$WATCHER_PID_FILE"
        fi
        /usr/share/eventcenter/watcher.sh &
        logger -t eventcenter "watcher started (realtime mode)"
    fi

    # Reload cron
    kill -HUP $(head -1 /var/run/crond.pid 2>/dev/null) 2>/dev/null

    logger -t eventcenter "service started (interval=${interval}min, realtime=${realtime:-off})"
}

stop_service() {
    # Stop watcher
    if [ -f "$WATCHER_PID_FILE" ]; then
        local _pid
        _pid=$(cat "$WATCHER_PID_FILE")
        if kill -0 "$_pid" 2>/dev/null; then
            kill -TERM "$_pid" 2>/dev/null
            logger -t eventcenter "watcher stopped (pid=$_pid)"
        fi
        rm -f "$WATCHER_PID_FILE"
    fi

    # Remove cron tasks
    grep -v 'eventcenter_check' /etc/crontabs/root > /tmp/ec_crs_$$ 2>/dev/null
    cp /tmp/ec_crs_$$ /etc/crontabs/root
    rm -f /tmp/ec_crs_$$

    # Reload cron
    kill -HUP $(head -1 /var/run/crond.pid 2>/dev/null) 2>/dev/null

    logger -t eventcenter "service stopped"
}

restart() {
    stop_service
    start_service
}
