# root/etc/init.d/eventcenter
#!/bin/sh /etc/rc.common

START=99
STOP=10

WATCHER_PID_FILE="/tmp/eventcenter_watcher.pid"

start() {
    local enable
    enable=$(uci -q get eventcenter.global.enable)
    [ "$enable" != "1" ] && return

    # --- Cron job setup ---
    setup_cron_jobs

    # --- Watcher process ---
    local realtime
    realtime=$(uci -q get eventcenter.openclash.realtime 2>/dev/null)
    if [ "$realtime" = "1" ]; then
        # Start watcher if not already running
        if ! pgrep -f "watcher.sh" >/dev/null 2>&1; then
            /usr/share/eventcenter/watcher.sh &
            echo $! > "$WATCHER_PID_FILE"
            logger -t eventcenter "watcher started (PID: $!)"
        else
            logger -t eventcenter "watcher already running"
        fi
    fi

    logger -t eventcenter "service started (realtime=${realtime:-off})"
}

setup_cron_jobs() {
    local _tmp_cron="/tmp/ec_cron_$$"
    crontab -l 2>/dev/null | grep -v "eventcenter" > "$_tmp_cron" 2>/dev/null || : > "$_tmp_cron"

    # OpenClash monitoring
    local oc_enable
    oc_enable=$(uci -q get eventcenter.openclash.enable 2>/dev/null)
    if [ "$oc_enable" = "1" ]; then
        local interval
        interval=$(uci -q get eventcenter.openclash.interval 2>/dev/null)
        [ -z "$interval" ] && interval=5
        echo "*/${interval} * * * * /usr/bin/eventcenter check openclash" >> "$_tmp_cron"
    fi

    # Node health monitoring
    local health_enable
    health_enable=$(uci -q get eventcenter.health.enable 2>/dev/null)
    if [ "$health_enable" = "1" ]; then
        local health_interval
        health_interval=$(uci -q get eventcenter.health.interval 2>/dev/null)
        [ -z "$health_interval" ] && health_interval=5
        echo "*/${health_interval} * * * * /usr/bin/eventcenter check node-health" >> "$_tmp_cron"
    fi

    # System health monitoring
    local sys_health_enable
    sys_health_enable=$(uci -q get eventcenter.system_health.enable 2>/dev/null)
    if [ "$sys_health_enable" = "1" ]; then
        local sys_interval
        sys_interval=$(uci -q get eventcenter.system_health.interval 2>/dev/null)
        [ -z "$sys_interval" ] && sys_interval=5
        echo "*/${sys_interval} * * * * /usr/bin/eventcenter check system-health" >> "$_tmp_cron"
    fi

    # Device monitoring
    local device_enable
    device_enable=$(uci -q get eventcenter.device_monitor.enable 2>/dev/null)
    if [ "$device_enable" = "1" ]; then
        local device_interval
        device_interval=$(uci -q get eventcenter.device_monitor.interval 2>/dev/null)
        [ -z "$device_interval" ] && device_interval=2
        echo "*/${device_interval} * * * * /usr/bin/eventcenter check device-monitor" >> "$_tmp_cron"
    fi

    # Subscription expiry monitoring
    local sub_enable
    sub_enable=$(uci -q get eventcenter.sub.enable 2>/dev/null)
    if [ "$sub_enable" = "1" ]; then
        local sub_interval
        sub_interval=$(uci -q get eventcenter.sub.check_interval 2>/dev/null)
        [ -z "$sub_interval" ] && sub_interval=6
        echo "*/${sub_interval} * * * * /usr/bin/eventcenter check sub" >> "$_tmp_cron"
    fi

    # Aggregate flush
    local agg_enable
    agg_enable=$(uci -q get eventcenter.@global[0].aggregate_enable 2>/dev/null)
    if [ "$agg_enable" = "1" ]; then
        echo "*/5 * * * * /usr/bin/eventcenter aggregate-flush" >> "$_tmp_cron"
    fi

    crontab "$_tmp_cron" 2>/dev/null
    rm -f "$_tmp_cron"
}

stop() {
    # Remove cron entries
    local _tmp_cron="/tmp/ec_cron_stop_$$"
    crontab -l 2>/dev/null | grep -v "eventcenter" > "$_tmp_cron" 2>/dev/null || : > "$_tmp_cron"
    crontab "$_tmp_cron" 2>/dev/null
    rm -f "$_tmp_cron"

    # Kill watcher
    if [ -f "$WATCHER_PID_FILE" ]; then
        kill $(cat "$WATCHER_PID_FILE") 2>/dev/null
        rm -f "$WATCHER_PID_FILE"
    fi
    pkill -f "watcher.sh" 2>/dev/null

    logger -t eventcenter "service stopped"
}

restart() {
    stop
    sleep 1
    start
}
